# seefood.service
[Unit]
Description=Seefood Python Service
After=network-online.target
Wants=network-online.target

[Service]
# Working directory of your project
WorkingDirectory=/home/seefood/seefood-py

# Use full path to uv (no shell, no PATH lookup)
ExecStart=/home/seefood/.local/bin/uv run python -m src.main dev

# Environment settings
Environment="PYTHONUNBUFFERED=1"
# Uncomment if you have environment variables in .env:
# EnvironmentFile=/home/seefood/seefood-py/.env

# Always restart the service, even on clean exit
Restart=always
RestartSec=3

# Run as the 'seefood' user and group
User=seefood
Group=seefood

# Logging (simple persistent logs)
StandardOutput=append:/var/log/seefood.log
StandardError=append:/var/log/seefood-error.log

# Allow some time for clean shutdown
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target